# install/inits/systemd/system/tyk-pump.service
[Unit]
Description=Tyk Pump
 
[Service]
Type=simple
User=root
Group=root
ExecStart=/opt/tyk-pump/tyk-pump
Restart=always

[Install]
WantedBy=multi-user.target
